# gwe.apps/th02/supervisor.conf
[program:th02]
command=/usr/local/bin/th02.py

