# gwe.apps/coh/supervisor.conf
[program:th02]
command=/usr/local/bin/coh.py

